# debian/xapian-bridge.service
[Service]
ExecStart=/usr/bin/xapian-bridge
