# debian/xapian-bridge.socket
[Socket]
ListenStream=3004

[Install]
WantedBy=sockets.target
