# data/xapian-bridge.service.in
[Service]
ExecStart=@bindir@/xapian-bridge
